# systemd/baksnapper@.timer
[Unit]
Description=Run baksnapper@%i each night

[Timer]
OnCalender=00:20
Persistent=true

[Install]
WantedBy=timers.target


