# systemd/baksnapper-ssh@.timer
[Unit]
Description=Run baksnapper-ssh@%i each night.

[Timer]
OnCalendar=00:20
Persistent=true

[Install]
WantedBy=timers.target


